# supervisor/dev/supervisor_celery_flower.ini
[program:查看器]
command = celery -A auxiliary flower --port=5566 --loglevel=DEBUG
directory = ../auxiliary
environment = PATH="venv/bin"
stdout_logfile = logs/celery.flower.log
stderr_logfile = logs/celery.flower.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 60
priority = 996
stdout_logfile_maxbytes = 200KB
stderr_logfile_maxbytes = 200KB
stdout_logfile_backups = 10
