# supervisor/dev/supervisor_celery_beat.ini
[program:定时器]
command = celery -A auxiliary beat --loglevel=DEBUG
directory = ../auxiliary
environment = PATH="venv/bin"
stdout_logfile = logs/celery.beat.log
stderr_logfile = logs/celery.beat.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 60
priority = 997
stdout_logfile_maxbytes = 200KB
stderr_logfile_maxbytes = 200KB
stdout_logfile_backups = 10
