# supervisor/dev/supervisor_celery_flower.ini
[program:celery-flower]
command=celery -A auxiliary flower --port=5566
directory=../auxiliary
environment=PATH="venv/bin"
stdout_logfile=logs/celery.flower.log
stderr_logfile=logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=996