# supervisor/dev/supervisor_celery_worker.ini
[program:执行器]
command = celery -A auxiliary worker --loglevel=DEBUG -P eventlet -c 50
directory = ../auxiliary
environment = PATH="venv/bin"
stdout_logfile = logs/celery.worker.log
stderr_logfile = logs/celery.worker.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 60
priority = 998
stdout_logfile_maxbytes = 200KB  ; 日志文件最大50MB
stdout_logfile_backups = 10     ; 保留10个旧日志文件
