# supervisor/dev/supervisor_celery_beat.ini
[program:celery-beat]
command=celery -A auxiliary beat -l debug
directory=../auxiliary
environment=PATH="venv/bin"
stdout_logfile=logs/celery.beat.log
stderr_logfile=logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=997