# supervisor/dev/supervisor_celery_worker.ini
[program:celery-worker]
command=celery -A auxiliary worker --loglevel=info -P eventlet -c 50
directory=../auxiliary
environment=PATH="venv/bin"
stdout_logfile=logs/celery.worker.log
stderr_logfile=logs/celery.worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=998
logfile_maxbytes=2MB
