# supervisor/product/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /ptools/logs/nginx.out.log
stderr_logfile = /ptools/logs/nginx.err.log
startsecs = 10
stopwaitsecs = 60
priority = 998
stdout_logfile_maxbytes = 1MB  ; 日志文件最大50MB
stdout_logfile_backups = 10     ; 保留10个旧日志文件
