# supervisor/product/redis.ini
[program:redis]
command = /usr/bin/redis-server --port 6379
autostart = true
autorestart = true
redirect_stderr = true
