# supervisor/product/supervisor_celery_beat.ini
[program:celery-beat]
command = celery -A auxiliary beat -l INFO
directory = /ptools
stdout_logfile = /ptools/logs/celery.beat.log
stderr_logfile = /ptools/logs/celery.beat.err.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 60
priority = 997
stdout_logfile_maxbytes = 1MB  ; 日志文件最大50MB
stdout_logfile_backups = 10     ; 保留10个旧日志文件
